# deploy/el34-assessor.service
[Unit]
Description=el34 Real Assessor (tw2 battle deterministic checks via docker exec)
After=docker.service
Requires=docker.service

[Service]
Type=simple
# el34 도커 호스트(TUBEWAR_REF_TARGET_IP, 예: .211)에서 root 로 실행 — docker exec 필요.
ExecStart=/usr/bin/python3 /opt/tw2/el34_assessor.py
Environment=ASSESS_PORT=9201
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
