# contents/battle-scenarios/compliance-w04.yaml
id: "compliance-w04"
category: "compliance"
title: "[컴플라이언스 W4] ISO 27001 Annex A — 구성 약점 악용 vs 기술 통제 점검"
difficulty: "medium"
time_limit: 2700
battle_type: "1v1"
description: |
  컴플라이언스 4주차 — **ISO 27001 Annex A 기술 통제**. 보안 인증(A.8.5)·암호 사용(A.8.24)·구성 관리(A.8.9)·
  정보 접근 제한(A.8.3)·기술적 취약성 관리(A.8.8)·네트워크 보안(A.8.20)을 점검한다. 공격자는 구성/암호화의
  약점을 악용하고, 감사자(너)는 Annex A 통제를 구현하고 충족을 증적으로 입증한다.

  ▸ 이야기
    공격자가 약한 구성을 노려 침투하고(태그 `co04r1`) 비인가 cron(`/etc/cron.d/co04cron`)과 계정 `compa04` 를
    심는다. 감사자는 ① 네트워크 보안 통제, ② 구성 변경 탐지, ③ Annex A 종합 점검·증적으로 대응한다.

  ▸ 모드
    • solo : 한 학생이 구성 약점 악용(자기 인프라 표적)과 통제 구현을 모두. 대상 = 본인 인프라.
    • duel : Red 가 상대(Blue) 인프라 구성을 악용, Blue 는 자기 인프라에 통제 구현.

  ▸ 공격자 위치(중요)
    이 주차의 Red 미션은 **망 외부 공격자 VM(192.168.0.202)** 기준이다. el34 는 공격자의 출처 IP를
    fw→ips→web(WAF) 전 계층에 그대로 보존하므로(HAProxy 제거·NAT 출처소실 없음), Red 채점은 **대상
    (el34) 인프라에 남은 탐지 흔적**(WAF audit / IPS / 호스트 계정·포트)과 **출처 IP(192.168.0.202)·
    고유 태그 상관**으로 한다. 반드시 **고유 태그**를 붙여 앰비언트 노이즈와 구분되게 한다.

  ▸ 이번 주차 전용 마커
    악용 태그 **`co04r1-<본인id>`**, 비인가 cron **`/etc/cron.d/co04cron`**(마커 `co04`), 계정 **`compa04`**,
    Suricata 네트워크 통제 **sid 9804011**, Wazuh 탐지 **id 104041**. 종료 시 cleanup.

red_missions:
  - order: 1
    assess_target: opponent
    target_vm: web
    points: 20
    instruction: |
      ## 🎯 RED-1 [구성 약점 악용] — 약한 암호화/구성 공략

      ### 상황
      약한 암호 알고리즘·구성 오류는 침투 경로다. 너는 태그 `co04r1-<본인id>` 로 침투해 구성 통제 공백을 입증한다.

      ### 할 일 (외부 공격자 VM 192.168.0.202)
      ```bash
      curl -s -A "co04r1-kim" -H "Host: admin.6v6.lab" \
        "http://192.168.0.161/files/read?id=co04r1-kim%27%20OR%201=1--+&Submit=Submit"
      ```

      ### 채점 방법
      대상 WAF audit(modsec)에 태그 `co04r1` 가 포함된 흔적이 남았는지로 채점.

      ### 합격 기준
      대상 WAF audit 에 `co04r1` 태그가 포함된 흔적이 남는다.
    hint: "태그(co04r1) 침투 → 대상 modsec audit. 채점: log_contains modsec co04r1."
    verify:
      type: "log_contains"
      expect: "co04r1"
      checks:
        - { id: "red-1-1", type: "log_contains", target: "web",
            params: { log: "modsec", pattern: "co04r1", since_sec: 1800 } }
      semantic:
        intent: "구성/암호 통제 미흡 악용(T1190) — Annex A 공백 입증. 태그(co04r1) 귀속."
        success_criteria:
          - "대상 WAF audit 에 co04r1 태그가 포함된 흔적"
          - "어떤 구성/암호 약점을 노렸는지 서술"
        acceptable_methods:
          - "외부 공격자 VM(192.168.0.202) 에서 태그 단 침투"
        negative_signs:
          - "태그 누락 / 대상 미도달"

  - order: 2
    assess_target: opponent
    target_vm: web
    points: 15
    instruction: |
      ## 🎯 RED-2 [구성 변조] — 비인가 cron + 계정

      ### 상황
      비인가 예약 작업(cron)은 ISO A.8.9 구성 관리 위반이다. 너는 비인가 cron **`/etc/cron.d/co04cron`**(마커 `co04`)과
      계정 **`compa04`** 를 심는다.

      ### 할 일
      ```bash
      echo '*/15 * * * * root /bin/bash -c "id # co04" >/dev/null 2>&1' | sudo tee /etc/cron.d/co04cron
      sudo useradd -m -s /bin/bash compa04
      ```

      ### 채점 방법
      대상 호스트에 `/etc/cron.d/co04cron`(마커 `co04`)이 있거나 `/etc/passwd` 에 `compa04` 가 있는지로 채점.

      ### 합격 기준
      `/etc/cron.d/co04cron`(마커 `co04`) 또는 계정 `compa04`.
    hint: "비인가 cron /etc/cron.d/co04cron (마커 co04) + useradd compa04. 채점: file_contains cron co04 + file_contains passwd compa04."
    verify:
      type: "file_contains"
      expect: "co04"
      checks:
        - { id: "red-2-1", type: "file_contains", target: "web",
            params: { path: "/etc/cron.d/co04cron", pattern: "co04" } }
        - { id: "red-2-2", type: "file_contains", target: "web",
            params: { path: "/etc/passwd", pattern: "compa04" } }
      semantic:
        intent: "구성 변조(T1053.003 cron + T1136 계정) — ISO A.8.9 구성 관리 위반. 비인가 예약 작업."
        success_criteria:
          - "/etc/cron.d/co04cron(마커 co04) 또는 계정 compa04"
          - "비인가 cron 의 위험과 A.8.9 연계 서술"
        acceptable_methods:
          - "cron.d 파일 + useradd"
        negative_signs:
          - "흔적 미생성"

blue_missions:
  - order: 1
    assess_target: self
    target_vm: ips
    points: 22
    instruction: |
      ## 🛡️ BLUE-1 [네트워크 보안 통제] — Suricata 룰 (sid 9804011) [A.8.20]

      ### 상황
      ISO A.8.20(네트워크 보안)을 구현한다. RED-1 침투를 잡는 Suricata 룰(**sid 9804011**)을 배포한다.

      ### 할 일
      ```
      alert http any any -> any any (msg:"COMPLIANCE-W04 ISO A.8.20 network intrusion"; flow:to_server; http.uri; content:"OR 1=1"; nocase; sid:9804011; rev:1;)
      ```
      ```bash
      sudo suricatasc -c reload-rules
      ```

      ### 채점 방법
      ips 의 `local.rules` 에 sid 9804011 룰이 있는지로 확인.

      ### 합격 기준
      `local.rules` 에 네트워크 보안 통제 sid 9804011 룰이 존재한다.
    hint: "네트워크 보안(A.8.20) Suricata 룰(sid 9804011) → reload. 채점: file_contains local.rules 9804011."
    verify:
      type: "file_contains"
      expect: "9804011"
      checks:
        - { id: "blue-1-1", type: "file_contains", target: "ips",
            params: { path: "/etc/suricata/rules/local.rules", pattern: "9804011" } }
      semantic:
        intent: "ISO A.8.20 네트워크 보안 — Suricata 침입 탐지 룰. 네트워크 통제 구현."
        success_criteria:
          - "local.rules 에 sid 9804011"
          - "A.8.20 네트워크 보안 통제 연계 언급"
        acceptable_methods:
          - "Suricata 룰 작성 + reload"
        negative_signs:
          - "패턴 부정확 / 문법 오류"

  - order: 2
    assess_target: self
    target_vm: siem
    points: 20
    instruction: |
      ## 🛡️ BLUE-2 [구성 변경 탐지] — Wazuh 룰 (id 104041) [A.8.9]

      ### 상황
      ISO A.8.9(구성 관리)를 구현한다. 비인가 cron·계정 변경을 잡는 Wazuh 룰(**id 104041**)을 만든다.

      ### 할 일
      ```xml
      <group name="edu,iso27001,">
        <rule id="104041" level="12">
          <if_group>edu</if_group>
          <match>cron.d</match>
          <description>COMPLIANCE-W04 - ISO A.8.9 config change: unauthorized cron/account</description>
        </rule>
      </group>
      ```
      ```bash
      sudo /var/ossec/bin/wazuh-control restart
      ```

      ### 채점 방법
      `local_rules.xml` 에 id 104041 룰이 있는지로 확인.

      ### 합격 기준
      `local_rules.xml` 에 구성 변경 탐지 id 104041 룰이 존재한다.
    hint: "구성 변경(A.8.9) Wazuh 룰(id 104041, cron.d FIM) → restart. 채점: file_contains local_rules.xml 104041."
    verify:
      type: "file_contains"
      expect: "104041"
      checks:
        - { id: "blue-2-1", type: "file_contains", target: "siem",
            params: { path: "/var/ossec/etc/rules/local_rules.xml", pattern: "104041" } }
      semantic:
        intent: "ISO A.8.9 구성 관리 — 비인가 cron/계정 변경 Wazuh FIM 룰. 구성 무결성 통제."
        success_criteria:
          - "local_rules.xml 에 rule id 104041"
          - "구성 변경 탐지 로직 + A.8.9 연계"
        acceptable_methods:
          - "Wazuh 룰 작성 + restart"
        negative_signs:
          - "탐지 로직 부실 / XML 오류"

  - order: 3
    assess_target: self
    target_vm: siem
    points: 18
    instruction: |
      ## 🛡️ BLUE-3 [Annex A 종합 점검·증적] — ISO 통제 충족 입증

      ### 상황
      ISO 27001 Annex A 기술 통제를 종합 점검하고 증적을 정리한다(SSH cipher/PAM, 커널 보안 파라미터, 디렉터리
      권한, 패치 현황, cron, 네트워크 설정, 부팅 서비스).

      ### 할 일
      보고서에:
      - 통제별 점검(`sshd -T ciphers`/`sysctl`/`ls -ld`/`apt list --upgradable`/`systemctl list-unit-files`) + Annex A ID 매핑 + 증적
      - 미흡 개선 + 비인가 cron(`/etc/cron.d/co04cron`)·계정(compa04) 헌팅
      ```bash
      sudo tail -200 /var/ossec/logs/alerts/alerts.json | \
        python3 -c 'import sys,json;[print(json.loads(l).get("rule",{}).get("id")) for l in sys.stdin if l.strip()]' | sort | uniq -c | sort -rn | head
      ```

      ### 채점 방법
      본인 SIEM 에 최근 경보가 수집돼 있는지 + 제출한 Annex A 종합 점검 증적으로 채점.

      ### 합격 기준
      본인 SIEM 에 최근 경보가 있다(+ Annex A 항목별 점검·증적·헌팅은 제출로 검증).
    hint: "ISO Annex A 점검(cipher/PAM/커널/권한/패치/cron/네트워크/부팅서비스) + 통제 ID 매핑 + 증적 + cron/계정 헌팅. 채점: wazuh_alert + 제출."
    verify:
      type: "wazuh_alert"
      expect: "alert"
      checks:
        - { id: "blue-3-1", type: "wazuh_alert", target: "siem", params: { since_sec: 1800 } }
      semantic:
        intent: "ISO 27001 종합 점검·증적 — Annex A 통제별 충족 입증 + 개선. 분석형 — 통제 ID/증적 제시 시 점수 상승."
        success_criteria:
          - "최근 경보 수집 확인"
          - "Annex A 통제별 점검(cipher/커널/권한/패치 등) + 통제 ID 매핑 + 증적"
          - "미흡 개선 + cron/계정 헌팅"
        acceptable_methods:
          - "점검 표 + 통제 ID + 증적 + 헌팅"
        negative_signs:
          - "통제 ID/증적 없는 서술 / 항목 누락"
