# scripts/tw2-tunnel.service
[Unit]
Description=tw2 Cloudflare Quick Tunnel (public -> UI :5173)
After=network-online.target tw2-ui.service
Wants=network-online.target
Requires=tw2-ui.service

[Service]
User=ccc
# quick tunnel: 계정/도메인 불필요, 매 기동 시 랜덤 https://*.trycloudflare.com 발급.
# 발급된 URL 은 stdout 에만 찍히므로 tee 로 .data/tunnel.log 에도 남겨 조회 가능하게 한다.
ExecStart=/bin/bash -c '/home/ccc/.local/bin/cloudflared tunnel --no-autoupdate --url http://localhost:5173 2>&1 | tee /home/ccc/tw2/.data/tunnel.log'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
