# scripts/tw2-tunnel.service
[Unit]
Description=tw2 Cloudflare Quick Tunnel (public -> UI :5173)
After=network-online.target tw2-ui.service
Wants=network-online.target tw2-ui.service
# 주의: Requires 로 tw2-ui 를 '하드' 의존하면, UI 재시작(빌드 반영) 시 터널까지 재기동돼
# quick tunnel 공개 URL 이 매번 새로 발급된다(→ 이전 URL 은 Cloudflare 1033). UI 와 느슨하게
# (Wants+After) 묶어 UI 를 재시작해도 터널·URL 은 그대로 유지되게 한다.

[Service]
User=ccc
# quick tunnel: 계정/도메인 불필요, 매 기동 시 랜덤 https://*.trycloudflare.com 발급.
# 발급된 URL 은 stdout 에만 찍히므로 tee 로 .data/tunnel.log 에도 남겨 조회 가능하게 한다.
ExecStart=/bin/bash -c '/home/ccc/.local/bin/cloudflared tunnel --no-autoupdate --url http://localhost:5173 2>&1 | tee /home/ccc/tw2/.data/tunnel.log'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
