# scripts/tw2-osd-tunnel.service
[Unit]
Description=tw2 OSD Cloudflare Quick Tunnel (public -> OpenSearch Dashboards :5602)
After=network-online.target
Wants=network-online.target

[Service]
User=ccc
ExecStart=/bin/bash -c '/home/ccc/.local/bin/cloudflared tunnel --no-autoupdate --url http://localhost:5602 2>&1 | tee /home/ccc/tw2/.data/osd-tunnel.log'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
